# provisioning/ansible/playbooks/attack/06_supplypro_ransom/templates/playbooks/lockscreen.service.j2
[Unit]
Description={{ attack_supplypro_ransom_payload_name }} starter


[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/{{ attack_supplypro_ransom_payload_name }}
User=%I

[Install]
WantedBy=multi-user.target
