# provisioning/ansible/playbooks/attack/00_phishing_sites/templates/_agent/ait_agent.service.j2
[Unit]
Description=AIT exercise logger

[Service]
Type=simple
ExecStart=/opt/ait_agent.py

[Install]
WantedBy=multi-user.target