# provisioning/ansible/playbooks/attack/ransom_lanturtle/files/attacker_0/auto-ncat.service
[Unit]
Description=Auto Ncat starter

[Service]
Restart=always
Restart=30s
ExecStart=/bin/bash /opt/ot_attack/ncat.sh

[Install]
WantedBy=multi-user.target
